# debian/maint/postinst
#! /bin/sh

set -e

if [ -x "/etc/init.d/gofish" ]; then
	update-rc.d gofish defaults >/dev/null
	/etc/init.d/gofish start
fi
